# haproxy-to-mondemand.crontab
* * * * * root /root/haproxy-to-mondemand-proxy --config_file /root/haproxy-to-mondemand.yaml
